# config/launchd/com.netlify-demo.quant.live-market-hours-test.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.netlify-demo.quant.live-market-hours-test</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/kenzhao/Projects/netlify-demo/scripts/run-live-market-hours-test.sh</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardErrorPath</key>
	<string>/Users/kenzhao/Projects/netlify-demo/docs/ai/logs/live-test.stderr.log</string>
	<key>StandardOutPath</key>
	<string>/Users/kenzhao/Projects/netlify-demo/docs/ai/logs/live-test.stdout.log</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Day</key>
		<integer>17</integer>
		<key>Hour</key>
		<integer>9</integer>
		<key>Minute</key>
		<integer>40</integer>
		<key>Month</key>
		<integer>6</integer>
	</dict>
</dict>
</plist>
